# install/consoled-dce.service
[Unit]
Description=Console Monitor DCE Service - Serial Proxy and Heartbeat Detection
Documentation=https://github.com/wiperi/consoled

# DCE 运行在 Console Server 端，依赖 database 服务
Requires=database.service
After=database.service config-setup.service
BindsTo=sonic.target
After=sonic.target

[Service]
Type=simple
User=root

# 使用统一入口，通过参数指定 DCE 模式
ExecStart=/usr/bin/consoled dce

# 自动重启配置
Restart=always
RestartSec=10

# 日志输出到 journal
StandardOutput=journal
StandardError=journal

# 需要 dialout 组权限访问串口
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
